# devel/.private/movo_network/etc/catkin/profile.d/50.movo_network_config.bash
. "/home/krish/catkin_ws/src/kinova-movo/movo_network/movo_network_config.bash"
