# devel/etc/catkin/profile.d/50.movo_config.bash
/home/krish/catkin_ws/devel/.private/movo_config/etc/catkin/profile.d/50.movo_config.bash